# onechat_project/config/onechat.conf
[program:onechat]
command=/root/projects/...
directory=/root/projects/onechat/onechat-project
user=onechat
autorestart=TRUE
redirect_stderr=true
stdout_logfile=/root/projects/onechat/logs/debug.log